# meta-yocto/recipes-bsp/beagleboard-audio/beagleboard-audio/beagleboard-audio
#!/bin/sh
### BEGIN INIT INFO
# Provides:          beagleboard mixer setup
# Required-Start:    $syslog
# Required-Stop:     $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Initialize the beagleboard audio mixer
# Description:       Unmute FRONT and set volume to ~70%.
### END INIT INFO

# Author: Darren Hart <dvhart@linux.intel.com>
# Based on /etc/init.d/skeleton

PATH=/sbin:/usr/sbin:/bin:/usr/bin
DESC="Audio mixer settings"
NAME=beagleboard-audio
AMIXER=`which amixer`
SCRIPTNAME=/etc/init.d/$NAME

# Exit if amixer is not installed
[ -x "$AMIXER" ] || exit 0

do_start() {
	# Enable the Headset (Audio Out)
	$AMIXER sset "Headset" 2 > /dev/null
	$AMIXER sset "HeadsetL Mixer AudioL1" on > /dev/null
	$AMIXER sset "HeadsetR Mixer AudioR1" on > /dev/null
}

case "$1" in
start)
	echo "$NAME: setting default mixer settings."
	do_start
	;;
stop)
	;;
*)
	echo "Usage: $SCRIPTNAME {start|stop}" >&2
	exit 3
	;;
esac

exit 0
